# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/data/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=5
loglevel=info
pidfile=/tmp/supervisord.pid
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
startretries=5
startsecs=3
stopwaitsecs=10
stdout_logfile=/data/logs/nginx_stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/data/logs/nginx_stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3

[program:backend]
command=node /app/backend/dist/server.js
directory=/app
user=veteran
autostart=true
autorestart=true
startretries=5
startsecs=5
stopwaitsecs=15
stdout_logfile=/data/logs/backend_stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/data/logs/backend_stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
environment=
    NODE_ENV="production",
    PORT="3000"
stopasgroup=true
killasgroup=true

[program:worker]
command=node /app/backend/dist/workers/index.js
directory=/app
user=veteran
autostart=true
autorestart=true
startretries=5
startsecs=5
stopwaitsecs=15
stdout_logfile=/data/logs/worker_stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/data/logs/worker_stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
environment=
    NODE_ENV="production"
stopasgroup=true
killasgroup=true

[eventlistener:processes]
command=/bin/sh -c "echo ready && while read line; do echo $line; done"
events=PROCESS_STATE
